# ns1-scan/crontab
*/30 * * * * python /task.py >> /logfile